# uploader/s6-overlay/s6-rc.d/api-uploader/run
#!/command/with-contenv /bin/bash
if [ "$API_UPLOADER_ENABLED" != "true" ]; then
    echo "API Uploader is disabled. Exiting."
    exit 0
fi

echo "Starting API Uploader monitoring $CONSUMPTION_DIR..."

upload_file() {
    local file_path="$1"
    local file_name=$(basename "$file_path")
    [[ "$file_name" == .* ]] && return
    echo "Uploading $file_name to $PAPERLESS_URL..."
    curl -s -L -f -H "Authorization: Token $PAPERLESS_TOKEN" -F "document=@$file_path" "${PAPERLESS_URL}api/documents/post_document/"
    if [ $? -eq 0 ]; then
        if [ "$API_UPLOADER_ON_SUCCESS" = "archive" ]; then
            mv "$file_path" "$ARCHIVE_DIR/$(date +%Y-%m-%d_%H-%M-%S)_$file_name"
        else
            rm "$file_path"
        fi
    fi
}

# Initial scan for existing files
echo "Performing initial scan of $CONSUMPTION_DIR..."
find "$CONSUMPTION_DIR" -maxdepth 1 -type f | while read -r file; do
    upload_file "$file"
done

# Monitor for new files using inotifywait
echo "Setting up inotify watches on $CONSUMPTION_DIR..."
exec inotifywait -m -e close_write -e moved_to --format '%w%f' "$CONSUMPTION_DIR" | while read -r file; do
    if [ -f "$file" ]; then
        # Small delay to ensure file is fully written/settled
        sleep 1
        upload_file "$file"
    fi
done
