# uploader/s6-overlay/s6-rc.d/wsdd/run
#!/command/with-contenv /bin/bash

echo "Starting WSDD..."

exec wsdd -v
