# uploader/s6-overlay/s6-rc.d/samba/run
#!/command/with-contenv /bin/bash
if [ "$SAMBA_ENABLED" != "true" ]; then
    exec sleep infinity
fi
echo "Starting Samba..."
exec smbd -F --no-process-group --log-stdout
