# uploader/s6-overlay/s6-rc.d/samba/run
#!/command/with-contenv /bin/bash

echo "Starting Samba..."

exec smbd -F --no-process-group --log-stdout
