# uploader/s6-overlay/s6-rc.d/config-samba/up
/etc/s6-overlay/scripts/config-samba.sh
