# uploader/s6-overlay/s6-rc.d/vsftpd/run
#!/command/with-contenv /bin/bash

echo "Starting vsftpd..."

exec vsftpd /etc/vsftpd/vsftpd.conf
