# uploader/s6-overlay/s6-rc.d/vsftpd/run
#!/command/with-contenv /bin/bash
if [ "$FTP_ENABLED" != "true" ]; then
    echo "FTP is disabled. Exiting."
    exit 0
fi
echo "Starting vsftpd..."
exec vsftpd /etc/vsftpd/vsftpd.conf
