# uploader/s6-overlay/s6-rc.d/config-vsftpd/up
/etc/s6-overlay/scripts/config-vsftpd.sh
